# openwrt_additional/rc.local

# Remount rootfs to enable sync
mount -o remount,rw,sync /

exit 0